# gps.service
[Unit]
Description=Pub: 8280 - GPS location; Sub: 8290 - RTCM corrections

[Service]
ExecStart=/usr/bin/python /home/pi/RoverGPS/gps_rover.py
StandardOutput=syslog
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
